# Sentinel_Detection_As_Code/detections/05_Persistence/Persistence via Registry Run Key Modification.yaml
title: Persistence via Registry Run Key Modification
id: dodea-sig-007-persistence-via-regi
status: experimental
description: 'Detects modification of common registry locations used for persistence
  such as Run, RunOnce, and similar keys.

  '
author: adam.ring
date: 2025-07-31 00:00:00+00:00
tags:
- attack.t1547.001
- attack.persistence
- ckc.installation
- data.source.DeviceRegistryEvents
logsource:
  product: windows
  category: registry_event
detection:
  selection:
    query: "DeviceRegistryEvents\n| where RegistryKey has_any (\n    @\"HKLM\\\\Software\\\
      \\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\",\n    @\"HKCU\\\\Software\\\
      \\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\",\n    @\"HKLM\\\\Software\\\
      \\Microsoft\\\\Windows\\\\CurrentVersion\\\\RunOnce\"\n  )\n| where ActionType\
      \ =~ \"RegistryValueSet\"\n| project TimeGenerated, DeviceName, RegistryKey,\
      \ RegistryValueName, RegistryValueData, InitiatingProcessFileName"
  condition: selection
falsepositives:
- Software installations or IT management tools
level: medium
tactics:
- Persistence
techniques:
- T1547.001
lifecycle: experimental
