# Sentinel_Detection_As_Code/detections/05_Persistence/Suspicious Scheduled Task Creation.yaml
title: Suspicious Scheduled Task Creation
id: dodea-sig-009-suspicious-scheduled-task-creation
status: experimental
description: >
  Detects creation of scheduled tasks using schtasks.exe or
  PowerShell, which can be used for persistence.
author: adam.ring
date: 2025-07-31 00:00:00+00:00
tags:
  - attack.t1053.005
  - attack.persistence
  - ckc.installation
  - data.source.DeviceProcessEvents
logsource:
  product: windows
  category: process_creation
detection:
  selection:
    query: |-
      DeviceProcessEvents
      | where FileName in~ ("schtasks.exe", "powershell.exe")
      | where ProcessCommandLine has_any ("Create", "Register-ScheduledTask")
      | project TimeGenerated, DeviceName, AccountName, ProcessCommandLine
  condition: selection
falsepositives:
  - Legitimate task creation by administrators or IT tools
level: high
tactics:
  - Persistence
techniques:
  - T1053.005
lifecycle: experimental
